# internal/server/modelrouter.service
[Unit]
Description=modelrouter - OpenAI and Anthropic compatible proxy server

[Service]
Type=simple
ExecStart=%h/.local/bin/modelrouter
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
