# internal/server/modelrouter.service
[Unit]
Description=ModelRouter - Ollama-compatible proxy server

[Service]
Type=simple
ExecStart=/usr/local/bin/ModelRouter
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
